# supervisor/conf.d/lighttpd.conf
[program:lighttpd]
command=/bin/bash -c "/usr/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf"
stdout_events_enabled=true
stderr_events_enabled=true